# user/nm-applet.service
[Unit]
Description=nm-applet
After=tint2.service

[Service]
ExecStart=/usr/sbin/nm-applet

[Install]
WantedBy=mystuff.target
